# capture/webdata/debian/cron.daily
#!/bin/sh
#
# cron script to load the database mydb with fund data found onthe internet
#
#
#
LOG=/var/log/weblog
date >> $LOG
/usr/bin/python /usr/share/pyshared/webdata/loadfunds.py >> $LOG
date >> $LOG